# recipes-bsp/pulseaudio/pulseaudio/10-local.rules
# Power up bluetooth when hci0 is discovered
ACTION=="add", KERNEL=="hci0", RUN+="/bin/hciconfig hci0 up"